# supervisord.conf
[supervisord]
nodaemon=true

[program:xinetd]
command=/usr/sbin/xinetd -dontfork
user=root
autorestart=true
priority=100

[program:apache2]
command=/etc/init.d/apache2 start
user=root
priority=200
